# roles/zfs_setup/files/syncoid-replicate.service
[Unit]
Description=Syncoid Replication to External ZFS Archive

[Service]
Type=oneshot
ExecStart=/usr/local/bin/syncoid-replicate.sh
