# roles/zfs_setup/files/zfs-archive-import.service
[Unit]
Description=Import ZFS Archive Pool and Resume Syncoid

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zfs-archive-import.sh
